# battery-logs/deep-blue-logs/battery-monitor.service
[Unit]
Description=Battery State Monitor
After=redis.service librescoot-battery.service
Requires=redis.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/data/battery-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target